# infra/vps/security/gernetix-security-alert-monitor.timer
[Unit]
Description=Run GerNetiX security alert scan every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
